# etc/systemd/system/android-mount.service.d/10-kumano.conf
[Service]
# Mount diag as /mnt/idd
Environment=EXTRA_FSTAB=/vendor/etc/idd.fstab
